# scripts/pairing-daemon/psmove-pairing.service
[Unit]
Description=PS Move Controller Pairing Daemon
Documentation=https://github.com/adangert/JoustMania
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/psmove-pairing-daemon.sh
Restart=always
RestartSec=10

# Environment
Environment=POLL_INTERVAL=10
Environment=DEBUG=0
# Common psmoveapi location - adjust if built elsewhere
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/joustmania/psmoveapi/build

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=psmove-pairing

[Install]
WantedBy=multi-user.target
