# conf/supervisor/conf.d/joust.conf
[program:joustmania]
command=/home/pi/JoustMania/joust.sh
user=root
directory=/home/pi/JoustMania/
autorestart=true
startretries=1000
