# conf/supervisor/conf.d/joust.conf
[program:joustmania]
command=/home/aaron/JoustMania/joust.sh
directory=/home/aaron/JoustMania/
autorestart=true
startretries=1000
stderr_logfile=/var/log/joustmania-err.log
stdout_logfile=/var/log/joustmania-out.log
