# scripts/setup/joustmania.service
[Unit]
Description=JoustMania Docker Compose Application
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pi/JoustMania
ExecStartPre=/usr/bin/docker compose down --remove-orphans
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300
User=pi
Group=pi

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
